# templates/zapbot-bridge.service
[Unit]
Description=Zapbot Webhook Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/env bun __ZAPBOT_DIR__/bin/zapbot-launch.ts --checkout __PROJECT_DIR__
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=append:/var/log/zapbot-bridge.log
StandardError=append:/var/log/zapbot-bridge.log

[Install]
WantedBy=multi-user.target
