# templates/zapbot-bridge.service
[Unit]
Description=Zapbot Webhook Bridge
After=network.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
# Load secrets from ~/.zapbot/config.json (canonical path since sbd#160 / zap#320).
# If a legacy .env still exists alongside config.json it is loaded as optional override.
ExecStartPre=/bin/bash -c 'install -dm 700 /run/zapbot-bridge && jq -r "\"ZAPBOT_WEBHOOK_SECRET=\" + .webhookSecret, \"ZAPBOT_API_KEY=\" + .apiKey" "$HOME/.zapbot/config.json" > /run/zapbot-bridge/env'
EnvironmentFile=/run/zapbot-bridge/env
EnvironmentFile=-__PROJECT_DIR__/.env
ExecStart=/usr/bin/env bun __ZAPBOT_DIR__/bin/webhook-bridge.ts
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5
StandardOutput=append:/var/log/zapbot-bridge.log
StandardError=append:/var/log/zapbot-bridge.log

[Install]
WantedBy=multi-user.target
